# canvas_worker.ini
[program:canvas_worker]
command=/home/digirati/natbakke/bin/python canvas_worker.py
directory=/home/digirati/natbakke/bin/
numprocs=4
process_name=canvas_worker_%(process_num)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile = /tmp/canvas_worker.log